# systemd/account-system.service
[Unit]
Description=JS.Gripe Account System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/account-system
Environment=NODE_ENV=production
Environment=ACCOUNT_HOST=127.0.0.1
Environment=ACCOUNT_PORT=9100
Environment=ACCOUNT_BASE_PATH=/api/v1/myaccount
Environment=ACCOUNT_ALLOWED_ORIGIN=https://account.js.gripe
Environment=ACCOUNT_DB_PATH=/opt/account-system/data/accounts.sqlite3
ExecStart=/usr/local/bin/npm start
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
